# etc/rc.local
#!/bin/bash
#
# /etc/rc.local: Local multi-user startup script.
#
export PYTHONPATH=/home/jonas/.syspath:$PYTHONPATH
chmod 777 /sys/class/backlight/acpi_video0/brightness
amixer -q set Master 100%
amixer -q set PCM 25%
export http_proxy=localhost:8118
