# etc/rc.local
#!/bin/bash
#
# /etc/rc.local: Local multi-user startup script.
#
export DISPLAY=:0
export PYTHONPATH=/home/jonas/.syspath:$PYTHONPATH
amixer -q set Front unmute
